# _deploy/systemd/hs-coming-soon-deploy.timer
[Unit]
Description=Check GitHub every 5 min and auto-deploy HuffmanStacks coming-soon if changed

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
